# borg.service
[Unit]
Description=Borg Autonomous Work Engine
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/borg
ExecStart=%h/borg/borg-rs/target/release/borg-server
Restart=always
RestartSec=5
Environment=HOME=%h
Environment=PATH=%h/.cargo/bin:%h/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/bin
EnvironmentFile=-%h/.config/borg/borg.env

[Install]
WantedBy=default.target
